# dns-update.service
[Unit]
Description=DNS Update Service
Documentation=https://github.com/floyd68/dns-update
After=network.target
Wants=network.target

[Service]
Type=simple
User=dns-updater
Group=dns-updater
WorkingDirectory=/opt/dns-update
Environment=PATH=/opt/dns-update/venv/bin
ExecStart=/opt/dns-update/venv/bin/python /opt/dns-update/app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/dns-update/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Environment variables (override these in /etc/dns-update/env)
Environment=HOSTED_ZONE_ID=
Environment=DOMAIN_NAME=
Environment=AWS_ACCESS_KEY_ID=
Environment=AWS_SECRET_ACCESS_KEY=
Environment=AWS_DEFAULT_REGION=us-east-1
Environment=ENABLE_IP_VALIDATION=true
Environment=ALLOWED_IPS=
Environment=ALLOWED_SUBNETS=
Environment=ENABLE_PASSWORD_AUTH=true
Environment=AUTH_PASSWORD=
Environment=FLASK_SECRET_KEY=
Environment=FLASK_HOST=0.0.0.0
Environment=FLASK_PORT=5000
Environment=FLASK_DEBUG=false
Environment=DNS_TTL=300
Environment=LOG_LEVEL=INFO
Environment=DNS_LOG_FILE=/opt/dns-update/logs/dns_updates.log

# Load environment variables from file
EnvironmentFile=/etc/dns-update/env

[Install]
WantedBy=multi-user.target 